# systemd/borg.timer
[Unit]
Description=Borg backup timer

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=default.target
